# deploy/com.claude-code-proxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.claude-code-proxy</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/richtobey/personal/claude-code-proxy/.venv/bin/claude-code-proxy</string>
	</array>
	<key>WorkingDirectory</key>
	<string>/Users/richtobey/personal/claude-code-proxy</string>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/tmp/claude-code-proxy-stdout.log</string>
	<key>StandardErrorPath</key>
	<string>/tmp/claude-code-proxy-stderr.log</string>
</dict>
</plist>
